# supervisor/conf.d/cron.ini
[program:cron]
command=/usr/sbin/crond -f
directory=/
autostart=true
autorestart=true
redirect_stderr=true
priority=1
stdout_logfile=/tmp/crond.log