# supervisor/conf.d/nginx.ini
[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
directory=/etc/nginx
autostart=true
autorestart=true
redirect_stderr=true
priority=2
stdout_logfile=/tmp/nginx.log