# supervisor/conf.d/trojan.ini
[program:trojan]
command=sh /trojan/lib/trojan.sh
directory=/trojan
autostart=true
autorestart=true
redirect_stderr=true
priority=3
stdout_logfile=/tmp/trojan.log
stopasgroup=true
killasgroup=true