# supervisor/conf.d/trojan.ini
[program:trojan]
command=/trojan/trojan-core/trojan -config {{CONFIG_PATH}}
directory=/trojan/trojan-core
autostart=true
autorestart=true
redirect_stderr=true
priority=3
stdout_logfile=/tmp/trojan.log